# data/packagekit-offline-update.service.in
[Unit]
Description=Update the operating system whilst offline

DefaultDependencies=no
Requires=sysinit.target dbus.socket
After=sysinit.target dbus.socket systemd-journald.socket system-update-pre.target
Before=shutdown.target system-update.target
# See packagekit.service
ConditionPathExists=!/run/ostree-booted
ConditionPathExists=!/var/lib/PackageKit/upgrade_lock

[Service]
Type=oneshot
ExecStartPre=/bin/touch /var/lib/PackageKit/disable-rpm-triggers
ExecStart=@libexecdir@/pk-offline-update
ExecStartPost=-/bin/sh @libexecdir@/pk-invoke-filetriggers.sh
ExecStartPost=-/bin/rm /var/lib/PackageKit/disable-rpm-triggers
ExecStartPost=@libexecdir@/pk-offline-update --shutdown
ExecStopPost=-/bin/rm /var/lib/PackageKit/disable-rpm-triggers

FailureAction=reboot
